# production/systemd/gs-backend.service
[Unit]
Description=General School Backend (FastAPI/Gunicorn)
After=network.target postgresql.service
Wants=postgresql.service
# 무한 재시작 폭주 방지 — 60초 내 5회 죽으면 재시작 중단.
# (StartLimit*는 반드시 [Unit] 섹션. [Service]에 두면 systemd가 무시함.)
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__INSTALL_DIR__/backend
EnvironmentFile=__INSTALL_DIR__/.env
# gunicorn -w (worker 수): (2×cores + 1) but capped for 16GB RAM laptop.
# 8코어/16GB 기준 9 worker × ~400MB = ~3.6GB → MemoryMax=4G와 일치.
# 더 늘리면 메모리·DB connection pool 부담. 줄이면 동시처리 부족.
# --timeout 300: nginx /api/ proxy_read_timeout 300s와 일치 (PDF/CSV 등 느린 endpoint 대비).
ExecStart=__INSTALL_DIR__/backend/venv/bin/gunicorn app.main:app \
    --workers 9 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 127.0.0.1:8002 \
    --timeout 300 \
    --graceful-timeout 30 \
    --max-requests 2000 \
    --max-requests-jitter 200 \
    --access-logfile - \
    --error-logfile -
Restart=always
RestartSec=5
# 메모리 누수 보호 (장시간 운영 안전망)
MemoryMax=4G

[Install]
WantedBy=multi-user.target
