# production/systemd/gs-hocuspocus.service
[Unit]
Description=General School Hocuspocus (Yjs WebSocket)
After=network.target gs-backend.service
Wants=gs-backend.service
# 무한 재시작 폭주 방지 — 60초 내 5회 죽으면 중단. (StartLimit*는 [Unit] 섹션 필수)
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__INSTALL_DIR__/backend-hocuspocus
EnvironmentFile=__INSTALL_DIR__/.env
Environment=NODE_ENV=production
Environment=HOCUSPOCUS_PORT=1234
ExecStart=/usr/bin/node __INSTALL_DIR__/backend-hocuspocus/dist/server.js
Restart=always
RestartSec=5
MemoryMax=2G

[Install]
WantedBy=multi-user.target
